# src/com.oshomadesse.bookssummary.pull.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.oshomadesse.bookssummary.autopull</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/seihoushouba/Library/Scripts/books-summary-auto-pull.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>10</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/seihoushouba/Library/Logs/BooksSummary/stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/seihoushouba/Library/Logs/BooksSummary/stderr.log</string>
</dict>
</plist>
